# agent/debian/wott-agent.timer
[Unit]
Description=WoTT Agent Timer

[Timer]
OnUnitActiveSec=1h
Unit=wott-agent.service

[Install]
WantedBy=multi-user.target
