# debian/wott-agent.timer
[Unit]
Description=Run wott-agent hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target