# debian/wott-agent-self-update.service
[Unit]
Description=wott-agent self update

[Service]
Type=oneshot
ExecStart=/opt/wott/bin/wott-agent-self-update.sh
RemainAfterExit=True
