# debian/wott-agent.service
[Unit]
Description=WoTT Agent

[Service]
Type=simple
ExecStart=/usr/bin/wott-agent
