# debian/wott-agent.service
[Unit]
Description=WoTT Agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/wott-agent daemon
Restart=always
RestartSec=1
StartLimitInterval=200
StartLimitBurst=5

[Install]
Alias=wott-agent
WantedBy=multi-user.target