# debian/wott-agent.timer
[Unit]
Description=Run wott-agent hourly

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
Unit=wott-agent.service

[Install]
WantedBy=default.target